# crontab
0 8 * * 1-5 ln -f /etc/hosts_work /etc/hosts
0 20 * * 1-5 ln -f /etc/hosts_free /etc/hosts
